# files/data/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -c $PHP_INI_DIR -y /usr/local/etc/php-fpm.conf

[program:apache]
command=/usr/sbin/httpd -D FOREGROUND
environment=APACHE_LOG_DIR=/var/log/apache2
redirect_stderr=true

[program:cron]
command = /usr/sbin/crond -f
stdout_logfile	= /var/log/cron.log
stderr_logfile	= /var/log/cron.log
autorestart=true
